# network/bind/init/services/bind/run
#! /bin/sh

exec named -c /etc/named/named.conf -g -u named
