# network/bind/init/services/bind/run
#! /bin/sh

trap "kill -TERM -1" KILL TERM INT
named -c /etc/named/named.conf -g -u named &
wait
