# mail/rspamd/init/services/rspamd/run
#! /bin/sh

trap "kill -TERM -1" KILL TERM INT
rspamd -f -u rspamd -g rspamd &
wait
