# mail/rspamd/init/services/rspamd/run
#! /bin/sh

exec rspamd -f -u rspamd -g rspamd
