# network/openvpn/init/services/openvpn/run
#! /bin/sh

trap "kill -TERM -1" KILL TERM INT
runsvdir /init/configs &
wait
