# network/openvpn/init/services/openvpn/run
#! /bin/sh

exec runsvdir /init/configs
