# rhodecode/init/services/community/run
#! /bin/sh

cd /rhodecode/services/community
./profile/bin/gunicorn -c gunicorn_config.py -u rhodecode -g rhodecode --paster rhodecode.ini
