# mail/postfix/init/services/postfix/run
#! /bin/sh

trap "kill -TERM -1" KILL TERM INT
postfix start-fg &
wait
