# web/nginx/init/services/nginx/run
#! /bin/sh

trap "kill -TERM -1" KILL TERM INT
nginx &
wait
