# rhodecode/init/services/vcsserver/run
#! /bin/sh

cd /rhodecode/services/vcsserver
./profile/bin/gunicorn -c gunicorn_config.py -u rhodecode --paster vcsserver.ini
