# mail/amavisd/init/services/amavisd/run
#! /bin/sh

exec /usr/local/bin/amavisd \
  -c /etc/amavisd.conf \
  -c /etc/amavisd/amavisd.conf \
  foreground
