# db/postgres/init/services/postgres/run
#! /bin/sh

exec su-exec postgres postgres -D /usr/local/pgsql/data
