# db/postgres/init/services/postgres/run
#! /bin/sh

trap "kill -TERM -1" KILL TERM INT
su-exec postgres postgres -D /usr/local/pgsql/data &
wait
