# mail/opendkim/init/services/opendkim/run
#! /bin/sh

trap "kill -TERM -1" KILL TERM INT
opendkim -f -x /etc/opendkim/opendkim.conf &
wait
