# mail/opendkim/init/services/opendkim/run
#! /bin/sh

exec opendkim -f -x /etc/opendkim/opendkim.conf
