# jira/init/services/jira/run
#! /bin/sh

cd /opt/jira
exec su-exec jira:jira ./bin/start-jira.sh -fg
